# bma-cms/keystone.conf
[supervisord]
nodaemon=true
loglevel=debug

[include]
files = /etc/supervisor/conf.d/*.conf

[program:keystone]
command = /usr/bin/node keystone.js
directory = /home/bma/dist
