# crlat_selenium/resources/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:vnc4server_go]
command=/start.sh
priority=500
autostart=true
autorestart=unexpected

#[program:go-agent]
#command=start_go.sh
#priority=500
#startsecs=2